# knxd/trunc/knxd-0.1/DEBIAN/postinst
#! /bin/sh -e

update-rc.d knxd defaults
/etc/init.d/knxd restart
